# bin/webtime.service

[Unit]
Description=Web service to print the current time of this VM for troubleshooting purposes.

[Service]
ExecStart=/usr/local/bin/time.py 1404
Restart=always

[Install]
WantedBy=multi-user.target



